# mo11y-streamlit.service
[Unit]
Description=Mo11y Streamlit Application
After=network.target mo11y-mcp-server.service

[Service]
Type=simple
User=dallas
Group=dallas
WorkingDirectory=/home/dallas/mo11y
Environment="PATH=/home/dallas/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/dallas/venv/bin/streamlit run /home/dallas/mo11y/app_enhanced.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target


