# mo11y-telegram-bot.service
[Unit]
Description=Mo11y Telegram Bot - Connect Alex to Telegram
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=dallas
Group=dallas
WorkingDirectory=/home/dallas/mo11y
Environment="PATH=/home/dallas/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="TELEGRAM_BOT_TOKEN=7911694889:AAFBBOgJsP2iTR5WsKwryieqdNILuzpptCU"
Environment="TELEGRAM_USER_ID=7358003141"
Environment="MO11Y_CONFIG_PATH=/home/dallas/mo11y/config.json"
ExecStart=/home/dallas/venv/bin/python /home/dallas/mo11y/telegram_bot_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
